# systemd/histomics-flask.service
[Unit]
Description=Histomics Flask API (Gunicorn)
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/histomics-prod/app

# Environment
Environment="PATH=/opt/histomics-prod/venv/bin"
EnvironmentFile=/opt/histomics-prod/config/.env.prod

# Apache/Nginx will proxy to this
ExecStart=/opt/histomics-prod/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 4 \
    --worker-class sync \
    --timeout 300 \
    --access-logfile /opt/histomics-prod/logs/flask/access.log \
    --error-logfile /opt/histomics-prod/logs/flask/error.log \
    --log-level info \
    --capture-output \
    app:app

# Restart behavior
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
Restart=on-failure
RestartSec=5s

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/histomics-prod/logs /data

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
